# systemd/dshield_enrich-ingest.service
[Unit]
Description=Cowrie LLM enrichment + session rollup (one-shot ingest pass)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=dshield_enrich
Group=dshield_enrich
WorkingDirectory=/opt/dshield_enrich
EnvironmentFile=/opt/dshield_enrich/.env
# Gate the run on local LLM reachability — a flaky GPU box fails the unit fast
# instead of burning a partial pass. Cloud reachability is checked inside
# enrich (ping → degrade to local-only on failure), so cloud outages do not
# block the run.
ExecStartPre=/opt/dshield_enrich/.venv/bin/python -m dshield_enrich.cli --config /opt/dshield_enrich/config/default.yaml healthcheck --scope llm
# Command-level enrichment (local LLM, optional cloud escalation).
# Required; failure here stops the chain.
ExecStart=/opt/dshield_enrich/.venv/bin/python -m dshield_enrich.cli --config /opt/dshield_enrich/config/default.yaml enrich
# Session rollup. Reads new closed sessions and joins with the
# enrichment docs just written above. `-` prefix → failure logged but does
# not fail the unit (sessions are recoverable on the next run).
ExecStart=-/opt/dshield_enrich/.venv/bin/python -m dshield_enrich.cli --config /opt/dshield_enrich/config/default.yaml rollup sessions
StandardOutput=journal
StandardError=journal
# Ollama can be slow on first call; allow a long run.
TimeoutStartSec=2h
Nice=10
# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/dshield_enrich
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
