# systemd/dshield_enrich-ingest.timer
[Unit]
Description=Run Cowrie ingest (enrich + session rollup) hourly
Requires=dshield_enrich-ingest.service

[Timer]
OnCalendar=hourly
RandomizedDelaySec=300
Persistent=true
Unit=dshield_enrich-ingest.service

[Install]
WantedBy=timers.target
