# systemd/dshield_vector_analysis_cluster.service
[Unit]
Description=Cowrie command clustering + novel escalation (one-shot pass)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=dshield_vector_analysis
Group=dshield_vector_analysis
WorkingDirectory=/opt/dshield_vector_analysis
EnvironmentFile=/opt/dshield_vector_analysis/.env
# cluster must succeed; escalate is non-fatal (exits non-zero when cloud.enabled=false)
ExecStart=/opt/dshield_vector_analysis/.venv/bin/python -m dshield_vector_analysis.cli --config /opt/dshield_vector_analysis/config/default.yaml cluster
ExecStart=-/opt/dshield_vector_analysis/.venv/bin/python -m dshield_vector_analysis.cli --config /opt/dshield_vector_analysis/config/default.yaml escalate
StandardOutput=journal
StandardError=journal
TimeoutStartSec=2h
Nice=10
# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/dshield_vector_analysis
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
