# systemd/dshield_vector_analysis_cluster.timer
[Unit]
Description=Run Cowrie clustering + novelty escalation every 6 hours
Requires=dshield_vector_analysis_cluster.service

[Timer]
OnCalendar=*-*-* 00,06,12,18:00:00
RandomizedDelaySec=300
Persistent=true
Unit=dshield_vector_analysis_cluster.service

[Install]
WantedBy=timers.target
