# systemd/dshield_enrich-analytics.service
[Unit]
Description=Cowrie analytics — clustering, novelty escalation, playbook naming, campaign mining
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=dshield_enrich
Group=dshield_enrich
WorkingDirectory=/opt/dshield_enrich
EnvironmentFile=/opt/dshield_enrich/.env
# Pre-gate on local LLM reachability — `name playbooks` needs the local
# LLM. If it's down the rest of the pipeline still runs and produces fresh
# clusters / campaigns; naming is the only step that degrades. Keep
# ExecStartPre as a soft warning rather than a hard fail by using `-` so
# the analytics pipeline runs even if the LLM is briefly unavailable.
ExecStartPre=-/opt/dshield_enrich/.venv/bin/python -m dshield_enrich.cli --config /opt/dshield_enrich/config/default.yaml healthcheck --scope llm

# --- Command-level clustering + cloud re-escalation -------------------------
# cluster must succeed — every downstream rollup job depends on its output.
ExecStart=/opt/dshield_enrich/.venv/bin/python -m dshield_enrich.cli --config /opt/dshield_enrich/config/default.yaml cluster commands
# escalate is optional (cloud may be disabled or the daily budget exhausted).
ExecStart=-/opt/dshield_enrich/.venv/bin/python -m dshield_enrich.cli --config /opt/dshield_enrich/config/default.yaml escalate

# --- Session-level clustering + playbook naming -----------------------------
# Each step uses `-` prefix so a transient failure in one step does not block
# the rest of the chain. Inputs to each step are produced by an earlier hourly
# run (rollup sessions in the ingest service) or by the previous step here.
ExecStart=-/opt/dshield_enrich/.venv/bin/python -m dshield_enrich.cli --config /opt/dshield_enrich/config/default.yaml cluster sessions
ExecStart=-/opt/dshield_enrich/.venv/bin/python -m dshield_enrich.cli --config /opt/dshield_enrich/config/default.yaml name playbooks

# --- IP-level rollup + clustering -------------------------------------------
# IP clusters are unnamed actor profiles; playbook membership for an IP is
# derived from its sessions at query time.
ExecStart=-/opt/dshield_enrich/.venv/bin/python -m dshield_enrich.cli --config /opt/dshield_enrich/config/default.yaml rollup ips
ExecStart=-/opt/dshield_enrich/.venv/bin/python -m dshield_enrich.cli --config /opt/dshield_enrich/config/default.yaml cluster ips

# --- Multi-session campaign mining ------------------------------------------
# Behaviour (FP-growth over playbook bags) + infrastructure (shared-artifact
# components). No LLM dependency; tolerates failure.
ExecStart=-/opt/dshield_enrich/.venv/bin/python -m dshield_enrich.cli --config /opt/dshield_enrich/config/default.yaml mine campaigns --kind all

StandardOutput=journal
StandardError=journal
# Longer timeout — 8 sequential steps can take a while when the corpus grows.
TimeoutStartSec=4h
Nice=10
# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/dshield_enrich
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
