# systemd/dshield_prism-analytics.timer
[Unit]
Description=Run Cowrie analytics (clustering + escalate + IP rollup + naming) every 6 hours
Requires=dshield_prism-analytics.service

[Timer]
OnCalendar=*-*-* 00,06,12,18:00:00
RandomizedDelaySec=300
Persistent=true
Unit=dshield_prism-analytics.service

[Install]
WantedBy=timers.target
